# ragapp/cli/cli/config/templates/systemd/ragapp-cli-web.target
[Unit]
After=network.target
Wants={{ cli_name }}-ragapp-web.service {{ cli_name }}-node-socketio.service

[Install]
WantedBy=multi-user.target
