# ragapp/cli/cli/config/templates/systemd/ragapp-cli-redis.target
[Unit]
After=network.target
Wants={{ cli_name }}-redis-cache.service {{ cli_name }}-redis-queue.service

[Install]
WantedBy=multi-user.target
